# pdns-sysrepo.service.in
[Unit]
Description=PowerDNS YANG configurator
Before=pdns.service

[Service]
Type=notify
ExecStart=@prefix@/@bindir@/pdns-sysrepo
User=pdns-sysrepo
Group=pdns-sysrepo

# Sandboxing
CapabilityBoundingSet= # None! We might need CAP_NET_BIND_SERVICE for prometheus in the future
LockPersonality=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true

ProtectSystem=full
ReadWritePaths=/etc/powerdns

RestrictAddressFamilies=
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=~ @clock @debug @module @mount @raw-io @reboot @swap @cpu-emulation @obsolete

[Install]
WantedBy=pdns.service
